# redis-sync.service
[Unit]
Description=Redis Database Sync Service
Documentation=https://github.com/joypciu/cache-api
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cache-api
EnvironmentFile=/var/www/cache-api/.env.sync
ExecStart=/usr/bin/python3 /var/www/cache-api/sync_redis.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redis-sync

[Install]
WantedBy=multi-user.target
