# cache-api/lutfur.service
[Unit]
Description=Lutfur Cache API Service
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/services/lutfur/cache-api
Environment="PATH=/home/ubuntu/services/lutfur/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="REDIS_HOST=localhost"
Environment="REDIS_PORT=6379"
Environment="CACHE_TTL=3600"
ExecStart=/home/ubuntu/services/lutfur/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
